# gateway/gateway.service
[Unit]
Description=Gateway service

[Service]
ExecStart=/home/user/ICS-lab/gateway/run.sh
WorkingDirectory=/home/user/ICS-lab/gateway

[Install]
WantedBy=multi-user.target